# deploy/budget-bot.service
[Unit]
Description=Budget Notebook - Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/telegram-expense-tracker
ExecStart=/home/pi/telegram-expense-tracker/.venv/bin/python main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
