# deploy/budget-bot.service
[Unit]
Description=Budget Notebook - Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=brunobarroso
WorkingDirectory=/home/brunobarroso/Desktop/telegram-expense-tracker
ExecStart=/home/brunobarroso/Desktop/telegram-expense-tracker/.venv/bin/python main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
